# cron/crontab
* * * * * /usr/bin/python3 ~/global/cron/ghistory.py >> /tmp/ghistory.log 2>&1