# cron/crontab
* * * * * /usr/bin/python3 $HOME/global/cron/ghistory.py >> /tmp/ghistory.log 2>&1